# server/note-server.service
[Unit]
Description=Note Server
After=network.target mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/note-server
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
